# trading-bridge/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1024x768x16
autorestart=true
priority=10

[program:x11vnc]
command=/usr/bin/x11vnc -display :0 -forever -nopw -quiet
autorestart=true
priority=20
depends_on=xvfb

[program:fastapi]
command=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app/src
autorestart=true
priority=30

